# platform/openwrt/custom/hciconf/files/final.init
#!/bin/sh /etc/rc.common

START=99

start() {
    curl http://localhost > /dev/null # Kick off input
	echo -ne '\x11' > /dev/ttyS0 # Tack unit ID to top
}
