# platform/openwrt/custom/hciconf/files/final.init
#!/bin/sh /etc/rc.common

START=99

start() {
        curl http://locxalhost > /dev/null # Kick off input
}
