# platform/openwrt/custom/hciconf/files/hciconf.init
#!/bin/sh /etc/rc.common

START=99

boot () {
	# This turns input and output buffering off, and streams input from the
	# footboard into a temporary file where the php workers can retrieve
	# as needed.
	while ! ip a s dev br-lan; do sleep 1; done
	cd /www/wsserver
	stdbuf -i0 -o0 -e0 cat /dev/ttyS0 | tee /tmp/ttyIn | stdbuf -i0 -o0 -e0 php-cli ttyws.php &
}
