# platform/openwrt/custom/hciconf/files/etc/init.d/setbaud
#!/bin/sh /etc/rc.common

START=22

start() {
        stty -F /dev/ttyS0 115200
}
